# systemd/bt-bridge.service
[Unit]
Description=Bluetooth LE to Classic Bridge Daemon
Documentation=https://gitea.hemna.com/admin/pi-bt-bridge
After=bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
ExecStartPre=/usr/bin/bluetoothctl pairable on
ExecStartPre=/usr/bin/bluetoothctl discoverable on
ExecStart=/home/waboring/pi-bt-bridge/.venv/bin/python3 -m src.main
WorkingDirectory=/home/waboring/pi-bt-bridge
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Run as root for Bluetooth access
User=root
Group=root

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
